# debian/postinst
#!/bin/bash
systemctl enable debiandde-first-config.service | true